# devservices/programs.conf
[program:devserver]
command=cargo run
autostart=false
